# Service/transit-worker.service
[Unit]
Description=TransitWorker

[Service]
Type=notify
ExecStart=/home/ubuntu/TransitWorker/TransitWorker

[Install]
WantedBy=multi-user.target
